# OpenHPC-v3/docker/image/bc/compute/etc/systemd/system/slurmd.service.d/override.conf
[Unit]
# Ensure hostname is set before starting slurmd
After=hostname-setup.service
Requires=hostname-setup.service
