# Moodle-Simple/docker/moodle/moodle4-simple/etc/supervisord.d/moodle.ini
[program:apache]
command=/usr/sbin/httpd -DFOREGROUND

[program:php-fpm]
command=/usr/sbin/php-fpm -F