# Moodle/docker/base/etc/systemd/system/docker.service.d/docker.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd --insecure-registry ${VCCCTR_IPADDR}:5000 --insecure-registry ${VCCCTR_IPADDR}:5001
EnvironmentFile=/etc/sysconfig/docker
