# CoursewareHub/playbooks/roles/systemd-worker-schedule2/templates/cwh-worker-node.service.j2
[Unit]
Description={{ service_desc }}

[Service]
ExecStart=docker compose up
WorkingDirectory={{ worker_schedule_dir }}
Type=oneshot

[Install]
WantedBy=multi-user.target
