# OpenHPC-v3/playbooks/roles/compute-schedule/templates/ohpc-compute-schedule.timer.j2
[Unit]
Description={{ service_desc }}

[Timer]
OnCalendar=*:0/{{ update_interval }}

[Install]
WantedBy=timers.target
