# CoursewareHub/playbooks/roles/worker-schedule/templates/cwh-worker-node.service.j2
[Unit]
Description={{ service_desc }}

[Service]
ExecStart={{ path_docker_cli }} compose up
WorkingDirectory={{ vcnode_ctl_dir }}
Type=oneshot

[Install]
WantedBy=multi-user.target
