# CoursewareHub/playbooks/roles/worker-schedule/templates/cwh-worker-node.timer.j2
[Unit]
Description={{ service_desc }}

[Timer]
OnCalendar=*:0/{{ exec_frequency }}

[Install]
WantedBy=timers.target

