# CoursewareHub-not-recommended/template/power-ctl/cwh-worker-node.timer.j2
[Unit]
Description=Update worker node power status of CoursewareHub

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target

