# Moodle-Simple/docker/base/etc/systemd/system/user_init_command.service
[Service]
ExecStart=/usr/local/bin/execute_user_init_command.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
